# debian/postinst
#!/bin/bash

# Copyright 2016, DELL EMC, Inc.

set -e
SERVICES=("on-http" "on-taskgraph" "on-dhcp-proxy" "on-tftp" "on-syslog" "on-web-ui")

#############################################
# Stop on-xxx 
# Use Rackhd service to manager the status of on-xxx services
#############################################
stop_component_services(){
    for i in ${SERVICES[@]}; do
        service ${i} stop
    done
}

stop_component_services

# Automatically added by dh_installinit
if [ -x "/etc/init.d/rackhd" ]; then
    update-rc.d rackhd defaults >/dev/null
    invoke-rc.d rackhd start
fi
# End automatically added section

